# agent-gateway.service.example
[Unit]
Description=Agent Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
WorkingDirectory=/opt/agent-workflows
EnvironmentFile=/opt/agent-workflows/.env
ExecStart=/opt/agent-workflows/bin/agent-gateway --port 18789
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

